# systemd/nestos-installer-pre.target
[Unit]
Description=Prepare for NestOS Installer Target

Requires=nestos-installer-disable-device-auto-activation.service
